# backend/deploy/ec2/karbg-backend.service
[Unit]
Description=Kar.bg Django ASGI (Daphne)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/avtoborsa/backend
EnvironmentFile=/opt/avtoborsa/backend/.env
Environment=DJANGO_SETTINGS_MODULE=backend.settings
ExecStart=/opt/avtoborsa/.venv/bin/daphne -b 127.0.0.1 -p 8000 backend.asgi:application
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
